# ops/systemd/ekobase-public-edge-check.timer
[Unit]
Description=Run Ekobase public edge health check every 5 minutes

[Timer]
OnBootSec=2m
OnUnitActiveSec=5m
AccuracySec=30s
Unit=ekobase-public-edge-check.service
Persistent=true

[Install]
WantedBy=timers.target
